# .docker/rootfs/etc/s6/vsftpd/run
#!/bin/sh
exec /usr/sbin/vsftpd /etc/vsftpd/vsftpd.conf
